# cron-dashboard/systemd/cron-dashboard-backend.service
[Unit]
Description=Cron Dashboard FastAPI Backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/guidda/.openclaw/workspace/cron-dashboard/backend
ExecStart=/usr/bin/python3 -m uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=2
EnvironmentFile=-/home/guidda/.openclaw/workspace/cron-dashboard/backend/.env

[Install]
WantedBy=multi-user.target
